# systemd-service-template.service
[Unit]
Description=

[Service]
WorkingDirectory=

User=
Group=

ExecStart=

[Install]
WantedBy=multi-user.target
